# supervisord/conf.d/sshd.conf
[program:sshd]
priority=2
user=root
command=/usr/sbin/sshd -D
autorestart=true